# deploy/client/retoonsClient.conf
[program:retoons_client]
command=serve -l 3000 -s /var/www/html
autostart=true
autorestart=true