# deploy/prod/server/retoonsServer.conf
[program:retoons_server]
command=node --require 'dotenv/config' /var/www/server/server.js
autostart=true
autorestart=true